# daemon/daemon
#!/bin/bash
### BEGIN INIT INFO
# Provides:          Text2SIP-Plugin
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: This file configures the Loxberry Text2SIP-Plugin.
# Description:       This file configures the Loxberry Text2SIP-Plugin.
### END INIT INFO

PATH="/sbin:/bin:/usr/sbin:/usr/bin"

. /lib/lsb/init-functions

### END INIT INFO

# Version v2018.2.12
# 11.02.2018 16:31:38

plugindir="text2sip"
logfile="REPLACELBPLOGDIR/Text2SIP.log"
plugconfigfile="REPLACELBPCONFIGDIR/Text2SIP.cfg"
flagfile="REPLACELBPCONFIGDIR/modify.me"
datafolder="REPLACELBPDATADIR"
touch $logfile
chown loxberry:loxberry $logfile
chmod 666 $logfile

if [ -r $flagfile ]
then
  echo `date +"%b  %e %H:%M:%S "`"modify.me found. Installing Text2SIP Plugin."
  echo `date +"%b  %e %H:%M:%S "`"modify.me found. Installing."                                                                2>&1 >>$logfile
  bash "apt-get -y install --reinstall man-db"  2>&1 >>$logfile
  echo `date +"%b  %e %H:%M:%S "`"Adapt rights / owner for files in REPLACELBPCONFIGDIR "   2>&1 >>$logfile
  chown loxberry:loxberry REPLACELBPCONFIGDIR/*                                                               2>&1 >>$logfile
  chmod 664 REPLACELBPCONFIGDIR/Text2SIP*                                                                     2>&1 >>$logfile
  cd $datafolder                                                                                                               2>&1 >>$logfile
  mkdir wav
  chown loxberry:loxberry wav                                                                                                  2>&1 >>$logfile
  chmod 755 wav                                                                                                                2>&1 >>$logfile
  cp lib*.so.* /usr/lib/                                                                                                    2>&1 >>$logfile
  if [ $? -eq 0 ]
  then
    echo `date +"%b  %e %H:%M:%S "`"Copy of libraries libopal libpt libspeexdsp libcapi20 and libodbc successfully completed." 2>&1 >>$logfile
  else
    echo `date +"%b  %e %H:%M:%S "`"ERROR: Copy of libraries libopal libpt libspeexdsp libcapi20 or libodbc failed. That is bad!"  2>&1 >>$logfile
  fi
  echo `date +"%b  %e %H:%M:%S "`"Try to install packages"                                                                     2>&1 >>$logfile
  bash dpkg -i libttspico-data_*                                                                                               2>&1 >>$logfile
  if [ $? -eq 0 ]
  then
    echo `date +"%b  %e %H:%M:%S "`"Package installation of libttspico-data successfully completed."                           2>&1 >>$logfile
  else
    echo `date +"%b  %e %H:%M:%S "`"ERROR: Package installation of libttspico-data failed. That is bad!"                       2>&1 >>$logfile
  fi
  bash dpkg -i libttspico0_*                                                                                                   2>&1 >>$logfile
  if [ $? -eq 0 ]
  then
    echo `date +"%b  %e %H:%M:%S "`"Package installation of libttspico0 successfully completed."                               2>&1 >>$logfile
  else
    echo `date +"%b  %e %H:%M:%S "`"ERROR: Package installation of libttspico0 failed. That is bad!"                           2>&1 >>$logfile
  fi
  sudo dpkg -i libttspico-utils_*                                                                                              2>&1 >>$logfile
  if [ $? -eq 0 ]
  then
    echo `date +"%b  %e %H:%M:%S "`"Package installation of libttspico-utils successfully completed."                          2>&1 >>$logfile
  else
    echo `date +"%b  %e %H:%M:%S "`"ERROR: Package installation of libttspico-utils failed. That is bad!"                      2>&1 >>$logfile
  fi
  echo `date +"%b  %e %H:%M:%S "`"Installation completed!"                                                                       2>&1 >>$logfile
  rm $flagfile																											        2>&1 >>$logfile
else
  echo `date +"%b  %e %H:%M:%S "`"modify.me not found. Plugin modification already done. Exit."
  echo `date +"%b  %e %H:%M:%S "`"Plugin modification already done."                                                           2>&1 >>$logfile
  chmod 755 REPLACELBPHTMLAUTHDIR/bin/sipcmd                                     2>&1 >>$logfile
fi

exit 0
